# pipeline/launchd/com.calgarytransit.archiver.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.calgarytransit.archiver</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/blessingudor/Desktop/development/yyc-bus-audit/pipeline/.venv/bin/python</string>
        <string>/Users/blessingudor/Desktop/development/yyc-bus-audit/pipeline/archive_rt.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/blessingudor/Desktop/development/yyc-bus-audit/pipeline</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/blessingudor/Desktop/development/yyc-bus-audit/data/rt/archiver.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/blessingudor/Desktop/development/yyc-bus-audit/data/rt/archiver.err.log</string>
</dict>
</plist>
